# crontabs/root
#!/bin/sh
0 * * * * 	/bin/sh /smhi-pmp3g-mqtt/smhi-pmp3g-mqtt.sh
